# systemd/spliceproxy.service
[Unit]
Description=Splice Proxy

[Service]
ExecStart=/usr/local/bin/spliceproxy -c /etc/spliceproxy.yaml

[Install]
WantedBy=multi-user.target